# origin-etc-init/zookeeper.conf
#info
env job=zookeeper
start on zk

stop on runlevel [016]

respawn
respawn limit 1 5

#console log
console output

pre-start script
  echo ~~ `date` >> /opt/logs/$job.log
end script


script
  . /etc/profile

  exec \
  $ZOOKEEPER_HOME/bin/zkServer.sh start-foreground \
  >> /opt/logs/$job.log 2>&1
end script


post-start script
  touch /opt/log/$job.touch
end script

post-stop script
  rm /opt/log/$job.touch
end script

